# backend/nyu-lims.service
[Unit]
Description=NYU LIMS Backend API
After=network.target postgresql.service

[Service]
Type=exec
User=nyu-lims
Group=nyu-lims
WorkingDirectory=/opt/nyu-lims/backend
Environment=PATH=/opt/nyu-lims/backend/venv/bin
ExecStart=/opt/nyu-lims/backend/venv/bin/gunicorn -c gunicorn.conf.py app.main:app
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=3

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/nyu-lims/backend/logs

[Install]
WantedBy=multi-user.target
